# com.user.pr-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.pr-notifier</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Update this path to where you cloned the repo -->
    <string>/Users/YOUR_USERNAME/pr-notifier/pr-notifier</string>
  </array>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/.local/share/pr-notifier/launchd-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/.local/share/pr-notifier/launchd-stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
